# src/petalinux/users/sean/profile.d/sean-settings.sh
#!/bin/bash

## This adds to what Lucas already does. 

export PATH=/home/root/scripts/:/home/root/scripts2/:/home/root/zynq/bin/:${PATH}   